# cameo.service
[Unit]
Description=Cameo MJPG stream service
After=multi-user.target

[Service]
Type=simple
User=pi
Environment=LD_LIBRARY_PATH=/opt/mjpg-streamer/
ExecStart=/home/pi/Cameo/start_server.sh
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target