# cameo.service
[Unit]
Description=Cameo service
After=multi-user.target

[Service]
Type=simple
User=pi
ExecStart=sudo /home/pi/Cameo/start_stream.sh ; sudo /usr/bin/node /home/pi/Cameo/cameo_server.js &
Restart=always
RestartSec=60

[Install]
WantedBy=multi-user.target